# deploy/supervisor/mqtt-bridge.conf
# bridge:mqtt-listen sudah punya backoff reconnect sendiri di dalam (lihat
# MqttBridgeListen::handle()) — autorestart=true di sini murni jaring
# pengaman kalau prosesnya crash total (fatal error PHP), bukan jalur
# reconnect utama untuk downtime broker Mosquitto.
[program:ctb-mqtt-bridge]
process_name=%(program_name)s
command=php /var/www/creative-trees-billing/artisan bridge:mqtt-listen
directory=/var/www/creative-trees-billing
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/creative-trees-billing/storage/logs/mqtt-bridge.log
