# deploy/supervisor/reverb.conf
# Naikkan `minfds` di /etc/supervisor/supervisord.conf (bagian [supervisord])
# ke 10000 supaya Supervisor bisa membuka cukup file descriptor untuk koneksi
# WebSocket Reverb — lihat dokumentasi resmi Reverb "Process Management".
[program:ctb-reverb]
process_name=%(program_name)s
command=php /var/www/creative-trees-billing/artisan reverb:start
directory=/var/www/creative-trees-billing
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/creative-trees-billing/storage/logs/reverb.log
