# deploy/supervisor/scheduler.conf
# `schedule:work`, bukan cron + `schedule:run` tiap menit — units:poll-state
# jalan tiap 30 detik (sub-minute task, lihat DECISIONS.md Fase 5), yang
# butuh proses scheduler yang terus hidup sepanjang menit, bukan dipicu ulang
# tiap 60 detik oleh cron.
[program:ctb-scheduler]
process_name=%(program_name)s
command=php /var/www/creative-trees-billing/artisan schedule:work
directory=/var/www/creative-trees-billing
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/creative-trees-billing/storage/logs/scheduler.log
